# crontab
*/1 * * * * curl --request GET 'http://rastishka-new.stage.by/testcrontab'
